# target/linux/ramips/mt7621/base-files/etc/init.d/set-irq-affinity
#!/bin/sh /etc/rc.common

START=99

start() {
    echo "Setting IRQ affinity for wired-only router..."
    
    [ -d "/proc/irq/21" ] && echo "2" > "/proc/irq/21/smp_affinity" 2>/dev/null
    [ -d "/proc/irq/24" ] && echo "4" > "/proc/irq/24/smp_affinity" 2>/dev/null
    [ -d "/proc/irq/23" ] && echo "8" > "/proc/irq/23/smp_affinity" 2>/dev/null
    
    for iface in eth0 eth1 br-lan; do
        if [ -d "/sys/class/net/$iface/queues/rx-0" ]; then
            echo 5 > /sys/class/net/$iface/queues/rx-0/rps_cpus 2>/dev/null
            echo "Set RPS for $iface"
        fi
    done
    
    modprobe mtkhnat 2>/dev/null
    
    echo 32768 > /proc/sys/net/netfilter/nf_conntrack_buckets 2>/dev/null
    echo 16384 > /proc/sys/net/netfilter/nf_conntrack_expect_max 2>/dev/null
    
    echo "IRQ affinity set for wired interfaces only"
}

stop() {
    echo "IRQ affinity service stopped"
}
