# target/linux/ramips/mt7621/base-files/etc/init.d/set-irq-affinity
#!/bin/sh /etc/rc.common

START=99

start() {

	echo "8" > "/proc/irq/23/smp_affinity"

	modprobe mtkhnat

	echo 32768 > /proc/sys/net/netfilter/nf_conntrack_buckets
	echo 16384 > /proc/sys/net/netfilter/nf_conntrack_expect_max

}
